# service/openrc/rpcd
#!/sbin/openrc-run

# Copyright 2021 (C) AlexCr4ckPentest
# OpenRC service unit for remote-power-control server

name="remote-power-control"
description="Simple CLI utility for remote power control over LAN"
command="/usr/local/bin/rpc"
command_args="--server"
